# deploy/systemd/zentrox-api.service
[Unit]
Description=Zentrox API
After=network.target mysql.service mariadb.service
Wants=mysql.service mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/zentrox/back
Environment=NODE_ENV=production
EnvironmentFile=-/var/www/zentrox/back/.env
ExecStartPre=/usr/bin/bash -c 'until mysqladmin ping -h127.0.0.1 --silent 2>/dev/null; do sleep 1; done'
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=3
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
